# alexandrite/root/etc/skel/.config/autostart/fcitx.desktop
[Desktop Entry]
Type=Application
Name=Fcitx
Comment=Start Input Method
Exec=fcitx
StartupNotify=false
NoDisplay=true
Terminal=false

